# pi/supervisor/flask.conf
[program:my_flask]
command=/usr/bin/uwsgi --ini /usr/local/share/opendoor/application.ini
directory=/usr/local/share/opendoor
user=pi
autostart=true
autorestart=true
