# pi/supervisor/subtitle.conf
[program:subtitle]
command=/usr/bin/python /bin/subtitle_downloader.py
autostart=true
autorestart=true
user=pi
